# supervisor-app.conf
[program:scrapy]
command = /usr/local/bin/scrapy crawl detail