# supervisor-app.conf
[program:scrapy]
command = cd /app
command = scrapy crawl detail